# services/x402-server/com.blueagent.x402server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.blueagent.x402server</string>

  <key>ProgramArguments</key>
  <array>
    <string>/usr/local/bin/node</string>
    <string>/Users/shun/projects/blue-agent/services/x402-server/dist/index.js</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/shun/projects/blue-agent/services/x402-server</string>

  <!-- Load .env variables -->
  <key>EnvironmentVariables</key>
  <dict>
    <key>NODE_ENV</key>
    <string>production</string>
    <key>PORT</key>
    <string>3002</string>
  </dict>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/Users/shun/projects/blue-agent/services/x402-server/logs/stdout.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/shun/projects/blue-agent/services/x402-server/logs/stderr.log</string>
</dict>
</plist>
